# ops/systemd/mnemes.service
[Unit]
Description=Mnemes memory authority service
After=network.target

[Service]
Type=simple
EnvironmentFile=%h/.config/mnemes/server.env
ExecStart=%h/.cargo/bin/mnemes-server ${MNEME_PORT} ${MNEME_DATA_DIR}
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/mnemes
LockPersonality=true
RestrictSUIDSGID=true

[Install]
WantedBy=default.target